# mac-studio-flux2/com.cstone.comfyui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.cstone.comfyui</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/jax/Documents/ComfyUI-v1/.venv/bin/python</string>
    <string>/opt/homebrew/Caskroom/comfy/0.9.1,260514248i6bhdm/ComfyUI.app/Contents/Resources/ComfyUI/main.py</string>
    <string>--base-directory</string><string>/Users/jax/Documents/ComfyUI-v1</string>
    <string>--front-end-root</string><string>/opt/homebrew/Caskroom/comfy/0.9.1,260514248i6bhdm/ComfyUI.app/Contents/Resources/ComfyUI/web_custom_versions/desktop_app</string>
    <string>--listen</string><string>0.0.0.0</string>
    <string>--port</string><string>8199</string>
    <string>--disable-auto-launch</string>
    <string>--dont-print-server</string>
    <string>--disable-metadata</string>
  </array>
  <key>WorkingDirectory</key><string>/opt/homebrew/Caskroom/comfy/0.9.1,260514248i6bhdm/ComfyUI.app/Contents/Resources/ComfyUI</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/jax/comfyui-service.log</string>
  <key>StandardErrorPath</key><string>/Users/jax/comfyui-service.log</string>
</dict>
</plist>
